# files/default/dcos.target
[Unit]
Requires=cloud-final.service
After=cloud-final.service
[Install]
WantedBy=multi-user.target
